# systemd/mnemora.timer
[Unit]
Description=Run Mnemora birthday check daily at 9 AM Pacific

[Timer]
OnCalendar=*-*-* 09:00:00 America/Los_Angeles
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
