# scripts/com.mnemora.birthday-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mnemora.birthday-daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>PROJECT_DIR_PLACEHOLDER/scripts/daily-check.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
    <key>TimeZone</key>
    <string>America/Los_Angeles</string>
  </dict>
  <key>StandardOutPath</key>
  <string>PROJECT_DIR_PLACEHOLDER/logs/launchd-daily.log</string>
  <key>StandardErrorPath</key>
  <string>PROJECT_DIR_PLACEHOLDER/logs/launchd-daily-error.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>

